# packages/scripts/cloud/admin/eliza-agent-router.service
[Unit]
Description=Eliza Agent Router (subdomain routing lookup for cloud agents)
Documentation=https://github.com/elizaOS/eliza
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=deploy
Group=deploy
WorkingDirectory=/opt/eliza
EnvironmentFile=/opt/eliza/cloud/.env.local
# Use the installed workspace binary directly to avoid `npx` startup overhead.
Environment=NODE_OPTIONS=--max-old-space-size=512
ExecStartPre=/opt/eliza/packages/scripts/cloud/admin/ensure-generated-keywords.sh
ExecStart=/opt/eliza/node_modules/.bin/tsx --tsconfig /opt/eliza/packages/cloud/shared/tsconfig.json /opt/eliza/packages/scripts/cloud/admin/daemons/agent-router.ts

Restart=always
RestartSec=5

MemoryMax=768M
MemoryHigh=512M

StandardOutput=journal
StandardError=journal
SyslogIdentifier=eliza-agent-router

NoNewPrivileges=yes
ProtectSystem=strict
PrivateTmp=true
ReadWritePaths=/opt/eliza
ReadOnlyPaths=/home/deploy/.ssh

TimeoutStopSec=10
KillMode=mixed
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target
